# ops/systemd/gemrouterfe.service
[Unit]
Description=GemRouterFE service
Wants=network-online.target xvfb.service x11vnc.service novnc.service
After=network-online.target xvfb.service x11vnc.service novnc.service

[Service]
Type=simple
WorkingDirectory=%h/GemRouterFE
EnvironmentFile=-%h/GemRouterFE/.env
Environment=NODE_ENV=production
ExecStart=%h/GemRouterFE/scripts/start-gemrouterfe.sh
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=20

[Install]
WantedBy=default.target
