# ops/systemd/gemrouter-nightly-restart.service
[Unit]
Description=Nightly restart of GemRouter

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart gemrouter.service
