# ops/systemd/gemrouter.service
[Unit]
Description=GemRouter
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=funboy
WorkingDirectory=/home/funboy/GemRouter
EnvironmentFile=/home/funboy/GemRouter/.env
Environment=HOME=/home/funboy
Environment=NODE_ENV=production
Environment=PNPM_HOME=/home/funboy/.local/share/pnpm
Environment=NVM_DIR=/home/funboy/.nvm
Environment=PATH=/home/funboy/.nvm/versions/node/v23.3.0/bin:/home/funboy/.local/share/pnpm:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/bin/bash -lc "pnpm build"
ExecStart=/bin/bash -lc "exec node dist/index.js"
Restart=always
RestartSec=10
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
