# ops/systemd/leakrouter-nightly-restart.service
[Unit]
Description=Nightly restart of LeakRouter

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart leakrouter.service
