# ops/systemd/leakrouter-nightly-restart.timer
[Unit]
Description=Timer: restart LeakRouter nightly at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=leakrouter-nightly-restart.service

[Install]
WantedBy=timers.target
