# ops/systemd/novnc.service
[Unit]
Description=GemRouterFE noVNC bridge
Wants=x11vnc.service
After=x11vnc.service

[Service]
Type=simple
WorkingDirectory=%h/GemRouterFE
EnvironmentFile=-%h/GemRouterFE/.env
ExecStart=%h/GemRouterFE/scripts/start-novnc.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
